# initial/cron/crontab
0 9 * * * root /jobs/job.sh >> /var/log/cron.log 2>&1
